# systemd/awccd.service
[Unit]
Description=AWCC Daemon
Documentation=https://github.com/tr1xem/AWCC
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/awccd
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
User=root
Group=root
WorkingDirectory=/var/run/awcc

# Create runtime directory
RuntimeDirectory=awcc
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
